# autostart/xcape.desktop
[Desktop Entry]
Name=Xcape
Exec=xcape
Type=Application
